# tools/bluetooth/patch-nvm-bdaddr.service
[Unit]
Description=Patch QCA Bluetooth NVM BDADDR
DefaultDependencies=no
After=local-fs.target
Before=bluetooth.service
ConditionPathExists=/usr/local/bin/patch-nvm-bdaddr.py
ConditionPathExists=/lib/firmware/qca/wcnhpnv21g.bin

[Service]
Type=oneshot
ExecStart=/usr/local/bin/patch-nvm-bdaddr.py

[Install]
WantedBy=multi-user.target
